# voltassistant-loadmanager/rootfs/etc/services.d/volt/finish
#!/usr/bin/execlineb -S1
# shellcheck shell=bash

/run/s6/basedir/bin/halt
